# systemd/homelab-run-wafflebot.timer
[Unit]
Description=Run WaffleBot daily at 4am server time

[Timer]
OnCalendar=04:00:00
Persistent=true

[Install]
WantedBy=timers.target 